# schedule_lab/echo.service
[Unit]
Description=Echo Service

[Service]
ExecStart=/usr/bin/nc -l 9090
